# supervisor.ini
[program:gourhunter]
user=yaping_pku
command=/home/yaping_pku/luxmusic/venv/bin/uwsgi --ini /home/yaping_pku/luxmusic/uwsgi.ini
directory=/home/yaping_pku/luxmusic/
autostart=true
autorestart=true
stderr_logfile=/mnt/sdb/log/supervisor.err.log
stdout_logfile=/var/sdb/log/supervisor.out.log
stopasgroup=true
