# supervisor.ini
[program:luxmusic]
user=yaping_pku
command=/home/yaping_pku/luxmusic/venv/bin/uwsgi --ini /home/yaping_pku/luxmusic/uwsgi.ini
directory=/home/yaping_pku/luxmusic/
autostart=true
autorestart=true
stderr_logfile=/mnt/sdb/log/luxmusic.err.log
stdout_logfile=/mnt/sdb/log/luxmusic.out.log
stopasgroup=true

[program:download]
user=yaping_pku
command=/home/yaping_pku/luxmusic/venv/bin/python /home/yaping_pku/luxmusic/manage.py download
directory=/home/yaping_pku/luxmusic/
autostart=true
autorestart=true
stderr_logfile=/mnt/sdb/log/download.err.log
stdout_logfile=/mnt/sdb/log/download.out.log
stopasgroup=true
